# .github/systemd/copilot-updater.timer
[Unit]
Description=Copilot — timer auto-update (toutes les heures)

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h
Unit=copilot-updater.service
Persistent=true

[Install]
WantedBy=timers.target
